# etc/cron.d/tfs-rescan
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
* * * * *   root    /home/torbox/torbox/lib/filesharing/tfs --scan
